# rpm/ngfd.service
[Unit]
Description=Non-Graphic Feedback Daemon
Requires=dbus.socket

[Service]
ExecStart=/usr/bin/ngfd
